# data/supervisor-confs/supervisor-smokeping.conf
[program:smokeping]
command=smokeping --nodaemon
