# data/supervisor-confs/supervisor-bind.conf
[program:named]
command=named -f -u bind
