# data/supervisor-confs/supervisor-dhcpd.conf
[program:dhcpd]
command=dhcpd -f -cf /etc/dhcp/dhcpd.conf -lf /etc/dhcp/dhcpd.leases
